# turkey-yacht-monitor/bot/aquila-watch-bot.service
[Unit]
Description=Aquila-Watch Telegram bot (long-polling)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aquila-watch-bot
EnvironmentFile=/opt/aquila-watch-bot/.env
ExecStart=/opt/aquila-watch-bot/venv/bin/python /opt/aquila-watch-bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
